# config/leader-key/com.machinespirit.leader-key.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the MachineSpirit Leader Key fork (the daily-driver launcher).

  RunAtLoad starts it at login; KeepAlive/SuccessfulExit=false makes it
  self-healing — launchd relaunches it if it ever CRASHES, but respects a
  deliberate Quit (clean exit) so you can stop it while developing the fork.

  __HOME__ is expanded to the real home path by install.sh (portability rule).
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.machinespirit.leader-key</string>
	<key>ProgramArguments</key>
	<array>
		<string>__HOME__/Applications/MachineSpirit Leader Key.app/Contents/MacOS/Leader Key</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
</dict>
</plist>
